# systemd/user/rclone-backup.service
[Unit]
Description=Rclone Backup Task

[Service]
ExecStart=/usr/bin/rclone copy -v /var/home/todd/文件 Koofr:documents
